# glogin.service
[Init]
Description=Login Prompt
After=gdm.target

[Service]
ExecStart=/usr/sbin/glogind
Restart=on-failure
